# ansible-redis-tls/roles/redis/files/redis.service.d/override.conf
# roles/redis/files/redis.service.d/override.conf
[Service]
LimitNOFILE={{ redis_nofile_limit }}
